# ansible/templates/tomcat.service.j2
# Ansible generated this file from {{ template_path }}; do not edit locally!
{# If this comment is present, this is the template file: disregard this warning #}

[Unit]
After=network.target
{% if inventory.get('tomcat').get('name') %}
Description={{ inventory.tomcat.name }}
{% endif %}

[Install]
WantedBy=multi-user.target

[Service]

AmbientCapabilities=CAP_NET_BIND_SERVICE
CacheDirectoryMode=750
Group=tomcat
LogsDirectoryMode=750
NoNewPrivileges=true
PrivateDevices=yes
PrivateTmp=yes
PrivateUsers=yes
ProtectClock=yes
ProtectControlGroups=yes
ProtectHome=read-only
ProtectHostname=yes
ProtectKernelLogs=yes
ProtectKernelTunables=yes
ProtectSystem=strict
Restart=always
RestartSec=10
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
SystemCallArchitectures=native
SystemCallFilter=@system-service
Type=forking
UMask=0007
User=tomcat

{% if inventory.get('tomcat').get('name') %}
CacheDirectory={{ inventory.tomcat.name }}
Environment="CATALINA_BASE=/var/lib/{{ inventory.tomcat.name }}"
Environment="CATALINA_HOME=/opt/{{ inventory.tomcat.name }}"
Environment="CATALINA_PID=/var/lib/{{ inventory.tomcat.name }}/temp/tomcat.pid"
{% if inventory.get('tomcat').get('jdk') %}
Environment="JRE_HOME=/usr/lib/jvm/{{ inventory.tomcat.jdk }}-ca-amd64/jre"
{% endif %}
ExecStart=/opt/{{ inventory.tomcat.name }}/bin/startup.sh
ExecStop=/opt/{{ inventory.tomcat.name }}/bin/shutdown.sh
LogsDirectory={{ inventory.tomcat.name }}
ReadOnlyPaths=/home/tomcat /opt/{{ inventory.tomcat.name }}
ReadWritePaths=/var/lib/{{ inventory.tomcat.name }} /var/log/{{ inventory.tomcat.name }}
RuntimeDirectory={{ inventory.tomcat.name }}
SyslogIdentifier={{ inventory.tomcat.name }}
WorkingDirectory=/var/lib/{{ inventory.tomcat.name }}
{% else %}
#
# This host has been configured, by Ansible, to be a db-only server.
# The following commands make `systemctl disable tomcat` return 0,
# while limiting resource leaks if we accidentally enable the service.
#
ExecStart=/bin/true
ExecStop=/bin/true
{% endif %}
